# debian/local/etc_init.d/glibc.sh
#! /bin/sh
#
### BEGIN INIT INFO
# Provides:          glibc
# Required-Start:
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: check for deprecated kernel versions
# Description:       This script detects deprecated kernel versions incompatible with
#                    the current version of the glibc
### END INIT INFO
#

set -e

# glibc kernel version check: KERNEL_VERSION_CHECK

: exit 0
